# podman-setup/container-metube.service
# container-a6aa7410fb20ecb8b57197248d44e34e57d9286d3de0118cd8cd1f903908979d.service
# autogenerated by Podman 5.3.1
# Thu Dec 12 22:12:47 EST 2024

[Unit]
Description=Podman container-a6aa7410fb20ecb8b57197248d44e34e57d9286d3de0118cd8cd1f903908979d.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start a6aa7410fb20ecb8b57197248d44e34e57d9286d3de0118cd8cd1f903908979d
ExecStop=/usr/bin/podman stop  \
	-t 10 a6aa7410fb20ecb8b57197248d44e34e57d9286d3de0118cd8cd1f903908979d
ExecStopPost=/usr/bin/podman stop  \
	-t 10 a6aa7410fb20ecb8b57197248d44e34e57d9286d3de0118cd8cd1f903908979d
PIDFile=/run/containers/storage/overlay-containers/a6aa7410fb20ecb8b57197248d44e34e57d9286d3de0118cd8cd1f903908979d/userdata/conmon.pid
Type=forking

[Install]
WantedBy=default.target
