# podman-setup/container-portainer_agent.service
# container-portainer_agent.service
# autogenerated by Podman 5.2.3
# Thu Oct 17 21:00:24 EDT 2024

[Unit]
Description=Podman container-portainer_agent.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=always
TimeoutStopSec=70
ExecStart=/usr/bin/podman run \
	--cidfile=%t/%n.ctr-id \
	--cgroups=no-conmon \
	--rm \
	--sdnotify=conmon \
	-d \
	-p 9001:9001 \
	--name portainer_agent \
	--replace \
	--security-opt label=disable \
	-v /run/podman/podman.sock:/var/run/docker.sock:Z portainer/agent:2.19.5
ExecStop=/usr/bin/podman stop \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
	-f \
	--ignore -t 10 \
	--cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
